# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[inet_http_server]
port=0.0.0.0:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://0.0.0.0:9001

[program:demo]
command=/home/vagrant/demo/venv/bin/gunicorn --chdir=/home/vagrant/demo/ -w 4 demo.wsgi:application -b 0.0.0.0:8000 --name demo --worker-class=gevent
directory=/home/vagrant/demo/
user=vagrant
autostart=true
startsecs=1
startretries=3 
stdout_logfile=stdout.log
stderr_logfile=stderr.log
